# docs/kmscon.service.in
[Unit]
Description=KMS System Console
Documentation=man:kmscon(1)

[Service]
ExecStart=@bindir@/kmscon --login -- /sbin/agetty -o '-p -- \\u' --noclear -- -

[Install]
WantedBy=multi-user.target
